# deploy/systemd/qaevaluate-backup.timer
[Unit]
Description=Run QA Evaluate SQLite backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=qaevaluate-backup.service

[Install]
WantedBy=timers.target
